# gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn --workers 3 --bind 127.0.0.1:5000 wsgi
directory=/deploy/app